# tcpdump.service
[Unit]
Description=TCPDump SystemD Daemon
[Service]
ExecStart=/usr/local/bin/tcpdump.sh
Restart=always
[Install]
WantedBy=multi-user.target
